# ansible/roles/swarmstack/files/etc/systemd/system/docker.service.d/http-proxy.conf
[Service]
Environment="HTTP_PROXY=proxy.example.com:80" "HTTPS_PROXY=proxy.example.com:443" "NO_PROXY=example.com,127.0.0.1"
